# debian/postinst
#!/bin/sh

chmod 755 /etc/init.d/rfm12piphp

update-rc.d rfm12piphp defaults

exit 0
